# crontab.cron
@reboot python3 face.py
00 22 * * * python3 zipFiles.py
00 23 * * * python3 uploadServer.py