# crontab.cron
#* * * * * echo "teste cron" >/dev/null 2>&1
@reboot /bin/sh /home/ubuntu/poc-face-detect-rasperbery/cronFace.sh > /dev/null 2>&1
30 * * * * /bin/sh /home/ubuntu/poc-face-detect-rasperbery/cronZipFiles.sh > /dev/null 2>&1
35 * * * * /bin/sh /home/ubuntu/poc-face-detect-rasperbery/cronUpload.sh > /dev/null 2>&1